# root/etc/systemd/system/owncloud-monitor@.service
[Service]
Type=simple
ExecStart=/usr/local/bin/owncloud-monitor %i
Restart=on-failure

[Install]
WantedBy=multi-user.target
